# supervisor/conf.d/ipsec.sv.conf
[program:ipsec]
command=/etc/init.d/ipsec restart
redirect_stderr=true