# supervisor/conf.d/pppd-dns.sv.conf
[program:pppd-dns]
command=/etc/init.d/pppd-dns start
redirect_stderr=true