# supervisor/conf.d/xl2tpd.sv.conf
[program:xl2tpd]
command=/etc/init.d/xl2tpd restart
redirect_stderr=true